# supervisor-proxy-services.conf
[group:proxy-app]
programs=pg-server,api-server,frontend,nginx,local-proxy

[program:pg-server]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
directory=/var/lib/postgresql
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=10
stderr_logfile=/var/log/supervisor/pg-server.err.log
stdout_logfile=/var/log/supervisor/pg-server.out.log
user=postgres

[program:api-server]
command=/root/.nvm/versions/node/v24.15.0/bin/npx tsx artifacts/api-server/src/index.ts
directory=/workspace
environment=DATABASE_URL="postgresql://postgres:postgres@localhost:5432/proxydb",PORT="8080",NODE_OPTIONS="--require /workspace/undici-proxy-preload.cjs"
autostart=true
autorestart=true
startsecs=3
startretries=8
priority=40
stderr_logfile=/var/log/supervisor/api-server.err.log
stdout_logfile=/var/log/supervisor/api-server.out.log
user=root

[program:frontend]
command=/root/.nvm/versions/node/v24.15.0/bin/npx vite --config vite.config.ts --host 0.0.0.0
directory=/workspace/artifacts/proxy-dashboard
environment=PORT="5173",BASE_PATH="/"
autostart=true
autorestart=true
startsecs=3
startretries=5
priority=40
stderr_logfile=/var/log/supervisor/frontend.err.log
stdout_logfile=/var/log/supervisor/frontend.out.log
user=root

[program:nginx]
command=/usr/sbin/nginx -c /workspace/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
startsecs=1
startretries=3
priority=30
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
user=root

[program:local-proxy]
command=/root/.nvm/versions/node/v24.15.0/bin/node /workspace/local-proxy.js
directory=/workspace
autostart=true
autorestart=true
startsecs=1
startretries=3
priority=20
stderr_logfile=/var/log/supervisor/local-proxy.err.log
stdout_logfile=/var/log/supervisor/local-proxy.out.log
user=root
